# geiger.service
[Unit]
Description=Geiger counter service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=alarm
ExecStart=/usr/bin/python /home/alarm/HAB-Experiment/logger.py 

[Install]
WantedBy=multi-user.target
