# geiger.service
[Unit]
Description=Geiger counter service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=root
WorkingDirectory=/home/alarm/HAB-Experiment 
ExecStart=/usr/bin/python /home/alarm/HAB-Experiment/logger.py 
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
